# templates/repmgrd.conf.j2
# Generated by Ansible
# All handmade changes can be rewrited!

console none

start on (local-filesystems and net-device-up IFACE!=lo)
stop on [!12345]

exec start-stop-daemon --start --user {{ postgresql_superuser }} --exec /usr/bin/repmgrd -- --verbose --monitoring-history --config-file={{ postgresql_repmgr_conf }}
